# systemd/telegram-bot-gastos-llm.service.template
[Unit]
Description=Telegram Expense Tracker Bot
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/telegram-bot-gastos-llm
Environment="PATH=/home/pi/telegram-bot-gastos-llm/.venv/bin"
Environment="LOG_DIR=/var/log/telegram-bot-gastos-llm"
ExecStart=/home/pi/telegram-bot-gastos-llm/.venv/bin/python run.py
Restart=always
RestartSec=10

LogsDirectory=telegram-bot-gastos-llm
LogsDirectoryMode=0755

[Install]
WantedBy=multi-user.target
